# meta-heipc-examples/recipes-examples/heipc-cpp-server/files/src/grpc_server.service
[Unit]
Description=GRPC Server

[Service]
ExecStart=/usr/bin/grpc_server
Restart=on-failure

[Install]
WantedBy=multi-user.target